# snapshot-backup@.service
[Unit]
Description=BTRFS Snapshot Backup

[Service]
Type=oneshot
ExecStart=snapshot-backup %I
