# snapshot@.service
[Unit]
Description=BTRFS Snapshot

[Service]
Type=oneshot
ExecStart=snapshot %I
