# snapshot-backup@.service
[Unit]
Description=BTRFS Snapshot Backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/snapshot-backup %I
