# snapshot@.service
[Unit]
Description=BTRFS Snapshot

[Service]
Type=oneshot
ExecStart=/usr/bin/snapshot %I || /usr/local/bin/snapshot %I
