# snapshot@.timer
[Unit]
Description=Run BTRFS snapshot

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
