# btrfssnapshottools/snapshot-backup@.timer
[Unit]
Description=BTRFS Snapshot Backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
